# launchd/local.translate-service.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.translate-service.worker</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>scripts/run_scheduled_worker.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__WORKDIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>UV_BIN</key>
      <string>__UV_BIN__</string>
      <key>WORKER_ACTIVE_START_TIME</key>
      <string>__WORKER_ACTIVE_START_TIME__</string>
      <key>WORKER_ACTIVE_END_TIME</key>
      <string>__WORKER_ACTIVE_END_TIME__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>0</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__WORKDIR__/logs/worker.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__WORKDIR__/logs/worker.stderr.log</string>
  </dict>
</plist>
